# supervisor_template.conf
[program:{project}]
command = /home/{project}/gunicorn_start
user = {project} ; Usuario que corre la aplicación
stdout_logfile = /home/{project}/logs/gunicorn_supervisor.log ; Donde se escribira la aplicación
redirect_stderr = true
autostart=true
autorestart=true